# backend/systemd/ecm-tunnel.service
[Unit]
Description=Cloudflare Tunnel for ECM
After=network.target

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate run ecm-tunnel

# 自動再起動設定
Restart=always
RestartSec=10

# ログ設定
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target